# charm/openvim/layer-openvim/templates/openvim.service
[Unit]
Description=openvim

[Service]
User=openvim
ExecStart=/opt/openmano/openvim/openvimd.py -c /opt/openmano/openvim/openvimd.cfg
Restart=always

[Install]
WantedBy=multi-user.target
